# test/deploy/linux/java/install/rhel/roles/configure/templates/supervisord.conf
[program:{{ service_id }}]
user=root
directory=/home/{{ ansible_user }}/
command=/opt/jdk-14.0.1/bin/java -jar appdynamics-cmdb-service.jar
process_name={{ service_id }}
user={{ ansible_user }}
autostart=true
autorestart=true
startsecs=10
stopwaitsecs=10
startretries=10
stopsignal=INT
redirect_stderr=true
stdout_logfile={{ stdout_path }}
stdout_logfile_maxbytes=10mb
stdout_logfile_backups=5
killasgroup=true
stopasgroup=true
environment=JAVA_HOME="/opt/jdk-14.0.1"
