# test/deploy/linux/tomcat/supervisord/shScript/install/linux/roles/configure/templates/supervisord.conf
[program:{{ service_id }}]
user=root
directory=/home/{{ ansible_user }}/{{ service_id }}
command=/opt/tomcat/latest/bin/catalina.sh run
process_name={{ service_id }}
autostart=true
autorestart=true
startsecs=10
stopwaitsecs=10
startretries=10
stopsignal=INT
redirect_stderr=true
stdout_logfile=/opt/tomcat/latest/logs/stdout.log
stdout_logfile_maxbytes=10mb
stdout_logfile_backups=5
killasgroup=true
stopasgroup=true
environment=JAVA_HOME="/opt/jdk-14.0.1",CATALINA_OUT=/dev/null,CATALINA_OPTS="-Xms300M -Xmx600M"
